# sandboxd/contrib/systemd/sandboxd.service
[Unit]
Description=sandboxd — isolated environment broker
Documentation=https://github.com/kontaktio/sandbox-daemon
After=docker.service
Wants=docker.service
# StartLimitIntervalSec belongs in [Unit] (systemd ≥ 230); placing it in
# [Service] generates a deprecation warning at unit-load time.
StartLimitIntervalSec=300

[Service]
Type=simple
User=sandbox
Group=sandbox

# Defense-in-depth for the unix socket mode: the daemon explicitly
# chmod()s sandboxd.sock to 0660 immediately after bind(), but
# UMask=0117 ensures the inode never even momentarily exists at a
# more permissive mode under a pre-exec race. The daemon sets this via the socket-path mechanism.
UMask=0117

# RuntimeDirectory creates /run/sandbox (tmpfs) with correct ownership/modes.
# The per-uid state directory /var/lib/sandboxd/<uid> is created explicitly
# by install.sh (StateDirectory does not fit the per-uid nested path model:
# it would chown the entire /var/lib/sandboxd chain to the sandbox user,
# breaking traversal by a co-resident sandbox-test daemon user).
RuntimeDirectory=sandbox
RuntimeDirectoryMode=0750

# Pin the daemon's ephemeral per-session runtime dirs (listener + event
# sockets) under the systemd RuntimeDirectory (/run/sandbox — tmpfs). Without
# these they default to $XDG_RUNTIME_DIR — UNSET for a system service — and
# silently fall back to $HOME (persistent disk), so ephemeral session
# artifacts would accumulate on disk across reboots. The socket and base-dir
# are pinned via --socket/--base-dir on ExecStart; these env vars pin the
# remaining XDG-derived runtime paths to the runtime tmpfs.
Environment=SANDBOX_LISTENER_DIR=/run/sandbox/listeners
Environment=SANDBOX_EVENTS_DIR=/run/sandbox/events

# Operators customize via `systemctl edit sandboxd` (drop-in override),
# not by editing this file. See § 4.3.
#
# NOTE: install.sh sed-substitutes @SANDBOX_BASE_DIR@ with the resolved
# per-uid path (/var/lib/sandboxd/<sandbox-uid>) at install time. The
# placeholder must never appear in the installed unit on a real host.
ExecStart=/usr/local/libexec/sandboxd/sandboxd \
    --base-dir @SANDBOX_BASE_DIR@ \
    --socket /run/sandbox/sandboxd.sock

# Restart policy: any crash bounces; stop bouncing if it crashes >5 times
# in 5 minutes (likely a config problem that won't fix itself).
Restart=on-failure
RestartSec=5s
StartLimitBurst=5

# Hardening — each directive has a one-line rationale in § 4.2.
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
DeviceAllow=/dev/kvm rw

[Install]
WantedBy=multi-user.target
