# sandboxd/contrib/systemd/sandboxd.service
[Unit]
Description=sandboxd — isolated environment broker
Documentation=https://github.com/kontaktio/sandbox-daemon
After=docker.service
Wants=docker.service
# StartLimitIntervalSec belongs in [Unit] (systemd ≥ 230); placing it in
# [Service] generates a deprecation warning at unit-load time.
StartLimitIntervalSec=300

[Service]
Type=simple
User=sandbox
Group=sandbox

# Defense-in-depth for the unix socket mode: the daemon explicitly
# chmod()s sandboxd.sock to 0660 immediately after bind(), but
# UMask=0007 ensures the inode never even momentarily exists at a
# more permissive mode under a pre-exec race. World bits are stripped
# (no-world), which is sufficient for the socket race-defense; the
# owner+group execute bits are preserved so the daemon can traverse
# the directories it creates (session dirs, lite tempdirs).
UMask=0007

# RuntimeDirectory creates /run/sandbox (tmpfs) with correct ownership/modes.
# The per-uid state directory /var/lib/sandboxd/<uid> is created explicitly
# by install.sh (StateDirectory does not fit the per-uid nested path model:
# it would chown the entire /var/lib/sandboxd chain to the sandbox user,
# breaking traversal by a co-resident sandbox-test daemon user).
RuntimeDirectory=sandbox
RuntimeDirectoryMode=0750

# Pin the daemon's ephemeral per-session runtime dirs (listener + event
# sockets) under the systemd RuntimeDirectory (/run/sandbox — tmpfs). Without
# these they default to $XDG_RUNTIME_DIR — UNSET for a system service — and
# silently fall back to $HOME (persistent disk), so ephemeral session
# artifacts would accumulate on disk across reboots. The socket and base-dir
# are pinned via --socket/--base-dir on ExecStart; these env vars pin the
# remaining XDG-derived runtime paths to the runtime tmpfs.
Environment=SANDBOX_LISTENER_DIR=/run/sandbox/listeners
Environment=SANDBOX_EVENTS_DIR=/run/sandbox/events

# Operators customize via `systemctl edit sandboxd` (drop-in override),
# not by editing this file. See § 4.3.
#
# NOTE: install.sh sed-substitutes @SANDBOX_BASE_DIR@ with the resolved
# per-uid path (/var/lib/sandboxd/<sandbox-uid>) at install time. The
# placeholder must never appear in the installed unit on a real host.
ExecStart=/usr/local/libexec/sandboxd/sandboxd \
    --base-dir @SANDBOX_BASE_DIR@ \
    --socket /run/sandbox/sandboxd.sock

# Restart policy: any crash bounces; stop bouncing if it crashes >5 times
# in 5 minutes (likely a config problem that won't fix itself).
Restart=on-failure
RestartSec=5s
StartLimitBurst=5

# Hardening — rationale for each directive below.
#
# NoNewPrivileges is intentionally absent: the daemon exec()s setcap helpers
# (sandbox-route-helper, sandbox-lima-helper) that carry file capabilities.
# Under NoNewPrivileges the kernel ignores file caps on execve, so the helpers
# run without their caps and every privileged operation fails (route-helper
# "setns: EPERM", lima-helper cannot setuid to the operator). The narrow
# setcap-helper model is the security boundary here, not NNP.
#
# ProtectHome is intentionally absent: it masks /home at the mount-namespace
# level for the entire service. The lima-helper setuid()s to the operator uid
# before exec'ing limactl, but stays in the service's mount namespace, so
# ~/.local/bin/limactl (Lima's default install path) becomes invisible.
# BindReadOnlyPaths cannot glob, so per-operator binds do not scale. The
# unprivileged daemon (running as the sandbox system user) already cannot read
# operator home directories via POSIX permissions, making ProtectHome redundant.
ProtectSystem=full
PrivateTmp=yes
DeviceAllow=/dev/kvm rw

[Install]
WantedBy=multi-user.target
